# luci/protocols/shadowvpn/root/etc/hotplug.d/iface/99-shadowvpn
#!/bin/sh

[ "$ACTION" = ifup ] || exit 0

shadowvpn_intf="`cat /var/etc/shadowvpn_intf`"

[ -n "$shadowvpn_intf" ] || exit 0
[ "$INTERFACE" != "$shadowvpn_intf" ] || exit 0

ifup $shadowvpn_intf
