# luci/applications/luci-shadowsocks/root/etc/hotplug.d/iface/99-shadowsocks
#!/bin/sh

[ "$ACTION" = ifup ] || exit 0

pgrep ss-redir >/dev/null 2>&1 && iptables -t nat -I zone_lan_prerouting -j SHADOWSOCKS