# luci/protocols/openconnect-ra/root/etc/hotplug.d/iface/99-openconnect
#!/bin/sh

[ "$ACTION" = "remove" ] || exit 0

PROTO="`ubus call network.interface.$INTERFACE status | grep proto | grep -q openconnect`"

[ $PROTO -eq 0 ] || exit 0

rm -f \
/etc/openconnect/ca-vpn-$INTERFACE.pem \
/etc/openconnect/user-cert-vpn-$INTERFACE.pem \
/etc/openconnect/user-key-vpn-$INTERFACE.pem
