# luci/luci-proto-tun2socks/root/etc/hotplug.d/iface/99-tun2socks
#!/bin/sh
TEMPFILE=/var/etc/tun2socks

[ "$INTERFACE" = wan ] || exit 0

[ -f "$TEMPFILE" ] || exit 0

[ "$ACTION" = ifup -o "$ACTION" = ifupdate ] && ifup "$(cat $TEMPFILE)"
