# luci/applications/luci-vsftpd-ra/root/etc/uci-defaults/luci-vsftpd
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@vsftpd-ra[-1]
	add ucitrack vsftpd-ra
	set ucitrack.@vsftpd-ra[-1].init=vsftpd-ra
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
