# luci/applications/luci-redsocks/root/etc/uci-defaults/luci-redsocks
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@redsocks[-1]
	add ucitrack redsocks
	set ucitrack.@redsocks[-1].init=redsocks
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
