# luci/applications/luci-redsocks2/root/etc/hotplug.d/iface/99-redsocks2
#!/bin/sh

[ "$ACTION" = ifup ] || exit 0

pgrep redsocks2 >/dev/null 2>&1 && iptables -t nat -I zone_lan_prerouting -j REDSOCKS2