# luci/applications/luci-aria2-ra/root/etc/uci-defaults/luci-aria2
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@aria2-ra[-1]
	add ucitrack aria2-ra
	set ucitrack.@aria2-ra[-1].init=aria2-ra
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
